# alpine_cnode_scripts_and_services/etc/init.d/node-exporter
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Node exporter service"

depend() {
	after network-online 
}

`grep -w 'USER' /etc/profile`
source /home/$USER/cnode_env
cd /home/$USER/node_exporter

start() {

        ebegin "Starting $RC_SVCNAME"
        start-stop-daemon --background --start --user $USER --exec ./node_exporter \
        --make-pidfile --pidfile /var/run/node-exporter.pid \
        -- --web.listen-address="$NODE_EXPORTER_LISTEN_ADDRESS" 
        eend $?
        
}

start_post() {
rcstatus=blank
	while [ -f /var/run/node-exporter.pid ]
	do
		sleep 10
		rcstatus=`rc-service $RC_SVCNAME status | awk '{print $NF}'`
		if [ -z $rcstatus ]; then
	                rc-service -c $RC_SVCNAME restart
			break
		fi
	done &
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
        start-stop-daemon --stop --user $USER --exec ./node_exporter \
        --pidfile /var/run/node-exporter.pid \
        -s 2 
        eend $?
}

