# alpine_cnode_scripts_and_services/etc/init.d/prometheus
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Prometheus service"

depend() {
	after network-online 
}

start() {

source /home/cardano/cnode_env

        ebegin "Starting $RC_SVCNAME"
        start-stop-daemon --background --start --exec /home/cardano/prometheus/prometheus \
        --make-pidfile --pidfile /var/run/prometheus.pid \
        -- --config.file="/home/cardano/prometheus/prometheus.yml" \
        --web.listen-address="$PROMETHEUS_LISTEN_ADDRESS" 
        eend $?
        
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
        start-stop-daemon --stop --exec /home/cardano/prometheus/prometheus \
        --pidfile /var/run/prometheus.pid \
        -s 2 
        eend $?
}

