# alpine_cnode_scripts_and_services/etc/init.d/prometheus
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Prometheus service"
command_user="cardano"
supervisor="s6"

depend() {
	after network-online 
}

`grep -w 'USER' /etc/profile`
source /home/$USER/cnode_env
cd /home/$USER/prometheus/

start() {
        ebegin "Starting $RC_SVCNAME"
        start-stop-daemon --start --background --user $USER --exec ./prometheus \
        --make-pidfile --pidfile /var/run/prometheus.pid \
        -- --config.file="/home/cardano/prometheus/prometheus.yml" \
        --web.listen-address="$PROMETHEUS_LISTEN_ADDRESS"
        eend $?
        
}

start_post() {
rcstatus=blank
        while [ -f /var/run/prometheus.pid ]
        do
                sleep 10
                rcstatus=`rc-service $RC_SVCNAME status | awk '{print $NF}'`
                if [ -z $rcstatus ]; then
                        rc-service -c $RC_SVCNAME restart
                        break
                fi
        done &
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
        start-stop-daemon --stop --user $USER --exec ./prometheus \
        --pidfile /var/run/prometheus.pid \
        -s 2 
        eend $?
}

